# contrib/systemd/spacenavd.service
[Unit]
Description=3Dconnexion Input Devices Userspace Driver

[Service]
Type=exec
PIDFile=/var/run/spnavd.pid
ExecStart=/usr/local/bin/spacenavd -d -l /dev/stdout

[Install]
WantedBy=graphical.target
